# scripts/com.odysseum.ticket-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Fires the check script every 15 minutes; the script's adaptive guard
     decides whether a new Pathé check is actually due (≈4 h baseline,
     every 15 min around an announced sale opening) and exits instantly
     otherwise, so most firings are ~2 s no-ops. Missed firings during
     sleep coalesce into one on wake; failed runs retry at the next firing.

     IMPORTANT: the job must run from a clone OUTSIDE ~/Documents (and other
     TCC-protected folders like Desktop/Downloads) — macOS blocks launchd
     agents there with "Operation not permitted". Setup:
       git clone https://github.com/<you>/odysseum-ticket-watch.git ~/.ticket-watch
       cp .env ~/.ticket-watch/ && cd ~/.ticket-watch
       python3 -m venv .venv && .venv/bin/pip install -r requirements.txt
       mkdir -p logs ~/Library/LaunchAgents
       cp scripts/com.odysseum.ticket-watch.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.odysseum.ticket-watch.plist
     Adjust the absolute paths below if your username/location differ. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.odysseum.ticket-watch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/serhii.afanasiev/.ticket-watch/scripts/local-check.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/serhii.afanasiev/.ticket-watch/logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/serhii.afanasiev/.ticket-watch/logs/launchd.log</string>
</dict>
</plist>
